# mutt/.local/share/systemd/user/mbsync@.service
[Unit]
Description=IMAP synchronization service

[Service]
type=oneshot
ExecStart=/usr/local/bin/mbsync %i:INBOX

[Install]
WantedBy=default.target
